# debian/alternc-slavedns.cron.d
# Synchronize the domain list for your local dns server
10 * * * * root /usr/sbin/alternc-slavedns > /dev/null
